# files/update-motd.d/00-greeting
#! /bin/sh
echo "------------------------------------------------------"
echo "--------------------- ISM DEVICE ---------------------"
echo "------------------------------------------------------"
